# ssh-monitor-watchdog.timer.example
[Unit]
Description=Run ssh-monitor watchdog every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
Unit=ssh-monitor-watchdog.service
Persistent=true

[Install]
WantedBy=timers.target
